# subserver/scripts/platform/linux/subserver.service
[Unit]
Description=subserver — 订阅分发服务
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__INSTALL_DIR__
ExecStart=/usr/bin/node __INSTALL_DIR__/index.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
