# smartxray/scripts/platform/linux/smartxray.service
[Unit]
Description=SmartXray - Xray Server Manager
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/xray run -c /var/lib/smartxray/config.json
StandardOutput=append:/var/log/smartxray/xray.log
StandardError=append:/var/log/smartxray/xray.log
RuntimeDirectory=smartxray
Restart=on-failure
RestartSec=5
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
[Unit]
Description=SmartXray - Xray Server Manager
After=network.target NetworkManager.service systemd-networkd.service

[Service]
Type=simple
User=__USER__
ExecStart=/usr/local/bin/xray run -c /home/__USER__/.config/smartxray/config.json
StandardOutput=append:/home/__USER__/.config/smartxray/logs/xray.log
StandardError=append:/home/__USER__/.config/smartxray/logs/xray.log
Restart=on-failure
RestartSec=5
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
